# poc-proxy.service
[Unit]
Description=PoC AI Proxy (SIP/RTP/TBCP walkie-talkie AI)
After=network.target asterisk.service
Wants=asterisk.service

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/Symbian-PTT-AI-chat
ExecStart=/usr/bin/python3 /home/ubuntu/Symbian-PTT-AI-chat/poc_proxy.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
TimeoutStartSec=10

[Install]
WantedBy=multi-user.target
