# poc-proxy.service
[Unit]
Description=PoC AI Proxy (SIP/RTP/TBCP walkie-talkie AI)
After=network.target asterisk.service
Wants=asterisk.service

[Service]
Type=simple
# INSTALL_DIR is written by install.sh into /etc/poc-proxy.env
EnvironmentFile=/etc/poc-proxy.env
WorkingDirectory=${INSTALL_DIR}
ExecStart=/usr/bin/python3 ${INSTALL_DIR}/poc_proxy.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
TimeoutStartSec=10

[Install]
WantedBy=multi-user.target
